# deploy/allinone/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid

# Postgres — local socket + 127.0.0.1 only
[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /data/pg
user=postgres
priority=10
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Neo4j — runs as the neo4j user, picks up /data/neo4j as its data path
[program:neo4j]
command=/usr/local/bin/run-neo4j.sh
user=neo4j
priority=15
autostart=true
autorestart=true
startsecs=15
environment=NEO4J_HOME="/var/lib/neo4j",JAVA_HOME="/usr/lib/jvm/java-17-openjdk-amd64"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ChromaDB — local-only HTTP server
[program:chromadb]
command=/opt/chroma/bin/chroma run --host 127.0.0.1 --port 8000 --path /data/chroma
user=root
priority=20
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# API server — needs Postgres + Neo4j healthy; entrypoint.sh waits before launching this
[program:api-server]
command=node --enable-source-maps /app/api/dist/index.mjs
user=root
priority=30
autostart=true
autorestart=true
startsecs=10
startretries=30
environment=PORT="3000",NODE_ENV="production"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# nginx — fronts everything to port 80
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
priority=40
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
